# nas-video-sync/nas-video-sync.timer
[Unit]
Description=Run NAS Video Sync every 24 hours
Requires=nas-video-sync.service

[Timer]
# Run 5 minutes after boot if the system was off during the scheduled time
OnBootSec=5min

# Run every 24 hours
OnUnitActiveSec=24h

# Alternative: Run at a specific time each day (uncomment and modify if preferred)
# OnCalendar=daily
# OnCalendar=*-*-* 02:00:00  # Run at 2 AM every day

# If the system was off during the scheduled time, run it as soon as possible
Persistent=true

[Install]
WantedBy=timers.target
