# nas-video-sync/nas-video-sync.service
[Unit]
Description=NAS Video Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Set working directory to where the script is located
WorkingDirectory=/opt/KCPAVid/VideoPi/nas-video-sync

# Source the configuration file and run the script
ExecStart=/bin/bash -c 'source /opt/KCPAVid/VideoPi/nas-video-sync/sync_nas_video.conf && /opt/KCPAVid/VideoPi/nas-video-sync/sync_nas_video.sh'

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nas-video-sync

# Security hardening (optional)
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
